# hamzero/etc/systemd/system/conbt.service
[Unit]
Description="Connect To BT Tether"

[Service]
Type=simple
ExecStart=/usr/bin/conbt

[Install]
WantedBy=multi-user.target
