# hamzero/DEBIAN/prerm
#!/usr/bin/bash
systemctl stop pat
systemctl disable pat
systemctl stop rigctld
systemctl disable rigctld
systemctl stop ardopc.path
systemctl disable ardopc.path
